# zero/meta-csp/recipes-cspd/systemd/files/cspd.service
[Unit]
Description=CSP Routing and execute Camera

[Service]
ExecStart=/usr/bin/csp_handler

[Install]
WantedBy=multi-user.target
